# mcp-hub.service
[Unit]
Description=MCP Servers Hub (FastAPI)
After=network.target

[Service]
Type=simple
User=forge
WorkingDirectory=/home/forge/mojodojo.guru
Environment=MCP_HOST=0.0.0.0
Environment=MCP_PORT=8000
ExecStart=/home/forge/mojodojo.guru/.venv/bin/python main.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target 